# syncd/scripts/persistent_mloop.conf
[program:persistent_mloop]
command=python3 /usr/bin/persistent_mloop.py
priority=50
stdout_logfile=NONE
stdout_syslog=true
stderr_logfile=NONE
stderr_syslog=true
autostart=true
autorestart=false
